# systemd/tcm-console.service
[Unit]
Description=TCM Console Manager
After=network.target

[Service]
Type=simple
User=tcm
Group=tcm
WorkingDirectory=/opt/tcm

ExecStart=/opt/tcm/venv/bin/python console/app.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tcm-console

Environment="PYTHONUNBUFFERED=1"
Environment="CONFIG_PATH=/etc/tcm/local-config.yaml"

[Install]
WantedBy=multi-user.target
